# deployment/supervisor/microblog.ini
[program:microblog]
command=/Users/king/workspace/microblog/venv/bin/gunicorn -b localhost:8000 -w 4 microblog:app
directory=/Users/king/workspace/microblog
user=king
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

stdout_logfile=/usr/local/var/log/microblog.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=30
stderr_logfile=/usr/local/var/log/microblog-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3